# wrestlebot/wrestlebot.service
[Unit]
Description=WrestleBot Data Collection Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/wrestlingdb/wrestlebot
Environment="PATH=/home/wrestlingdb/venv/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/home/wrestlingdb/wrestlebot/.env
ExecStart=/home/wrestlingdb/venv/bin/python /home/wrestlingdb/wrestlebot/main.py

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=0

# Resource limits
LimitNOFILE=65536
MemoryLimit=2G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wrestlebot

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
